# telegram-bot.service
[Unit]
Description=Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/telegram-bot
EnvironmentFile=-%h/telegram-bot/.env
ExecStart=%h/telegram-bot/venv/bin/python %h/telegram-bot/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
